# rootfs/etc/s6-overlay/s6-rc.d/chronyd-exporter/run
#!/bin/sh
exec s6-setuidgid chrony \
  /bin/chrony_exporter \
    --web.listen-address=:9123 \
    --chrony.address=unix:///run/chrony/chronyd.sock
